# systemd/update-dmenu-extended-db.timer
[Unit]
Description=Run update-dmenu-extended-db.service every 5 minutes

[Timer]
OnCalendar=*:0/20

[Install]
WantedBy=timers.target
